# centos7_lamp_stack_php7/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock 

[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D


[program:redis]
command=/usr/bin/redis-server /etc/redis.conf
autostart=true
autorestart=true
user=root
stdout_logfile=/var/log/redis/stdout.log
stderr_logfile=/var/log/redis/stderr.log


[program:mariadb]
command=/usr/bin/mysqld_safe

[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND

[group:allservices]
programs=sshd,mariadb,httpd,redis

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock         ; use a unix:// URL  for a unix socket

